# woof-code/rootfs-packages/puppycore_noarch/etc/rc.d/rc.local
#this file called from rc.sysinit
#you can edit this file
#When firewall is installed, will append lines to this file...
